# templates/systemd/lora-pkt-fwd.service
# Template LoRaCore — Packet Forwarder (sx1302_hal via USB)
# Fonte: docs/DOC_PROTOCOLO_COMUNICACAO_LORAWAN.md, Secao 4.4
# Destino: /etc/systemd/system/lora-pkt-fwd.service
#
# Substituir:
#   <USER>          — usuario do sistema (ex: marlon)
#   <PKT_FWD_PATH>  — diretorio do packet forwarder (ex: /home/marlon/packet_forwarder)

[Unit]
Description=LoRa Packet Forwarder (RAK2287 USB)
After=local-fs.target

[Service]
Type=simple
StandardOutput=journal
StandardError=journal
WorkingDirectory=<PKT_FWD_PATH>
ExecStartPre=/bin/sh -c "test -e /dev/ttyACM0"
ExecStart=<PKT_FWD_PATH>/lora_pkt_fwd
Restart=always
RestartSec=5
WatchdogSec=120
CPUWeight=150
MemoryHigh=64M
MemoryMax=128M
User=<USER>

[Install]
WantedBy=multi-user.target
